# systemd/jukebox.service
[Unit]
Description=Jukebox controller service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/jukebox
EnvironmentFile=/etc/jukebox/jukebox.env
ExecStart=/opt/jukebox/.venv/bin/python -m jukebox
Restart=on-failure
RestartSec=5
RestartPreventExitStatus=2
SupplementaryGroups=input
StandardInput=null
ReadWritePaths=/var/lib/jukebox

[Install]
WantedBy=multi-user.target
